# packages/daypass/files/uci-defaults/init
#!/bin/sh
# First-boot: generate a random mihomo API secret if none is set.
if [ -z "$(uci -q get __PKG_NAME__.main.api_secret)" ]; then
	uci -q set "__PKG_NAME__.main.api_secret=$(head -c 16 /dev/urandom | md5sum | cut -c1-16)"
	uci -q commit __PKG_NAME__
fi
exit 0
